# agent-runtime/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
loglevel=info

; ── 1. Virtual display ────────────────────────────────────────────────────────
[program:xvfb]
command=Xvfb :99 -screen 0 1280x800x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/xvfb.log
stderr_logfile=/data/logs/xvfb.log
startsecs=1

; ── 2. Window manager ────────────────────────────────────────────────────────
[program:openbox]
command=openbox --display :99
autostart=true
autorestart=true
priority=20
environment=DISPLAY=":99",HOME="/root"
stdout_logfile=/data/logs/openbox.log
stderr_logfile=/data/logs/openbox.log
startsecs=2

; ── 3. VNC server (captures Xvfb display) ────────────────────────────────────
[program:x11vnc]
; Uses /data/.vncpass written by start.sh from OPENCLAW_GATEWAY_TOKEN
command=x11vnc -display :99 -rfbauth /data/.vncpass -rfbport 5900 -listen localhost -xkb -ncache 10 -forever -shared -quiet -noxdamage
autostart=true
autorestart=true
priority=30
stdout_logfile=/data/logs/x11vnc.log
stderr_logfile=/data/logs/x11vnc.log
startsecs=3

; ── 4. noVNC WebSocket bridge ────────────────────────────────────────────────
[program:novnc]
command=websockify --web=/usr/share/novnc 6080 localhost:5900
autostart=true
autorestart=true
priority=40
stdout_logfile=/data/logs/novnc.log
stderr_logfile=/data/logs/novnc.log
startsecs=2

; ── 5. Nginx reverse proxy ────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=50
stdout_logfile=/data/logs/nginx.log
stderr_logfile=/data/logs/nginx.log

; ── 6. Freemi agent runtime ───────────────────────────────────────────────────
[program:agent]
command=node /app/dist/index.js
directory=/app
autostart=true
autorestart=true
priority=60
environment=DISPLAY=":99",HOME="/root",NODE_ENV="production",PLAYWRIGHT_BROWSERS_PATH="/ms-playwright"
stdout_logfile=/data/logs/runtime.log
stderr_logfile=/data/logs/runtime.log
startsecs=5
